# docker/config/supervisord-worker-staging.conf
[unix_http_server]
file=/var/www/catima/tmp/supervisor.sock

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/www/catima/tmp/supervisor.sock

[program:sidekiq]
command=bundle exec sidekiq -C /var/www/catima/config/sidekiq.yml
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/catima/log/sidekiq.log
stdout_logfile_maxbytes=0
