# docker/config/supervisord-app-staging.conf
[unix_http_server]
file=/var/www/catima/tmp/supervisor.sock

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/www/catima/tmp/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -c /var/www/catima/config/nginx.conf
user=catima
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:puma]
autorestart=true
command=bundle exec puma -C config/puma.rb -e staging
user=catima
stdout_logfile=/var/www/catima/log/puma.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/catima/log/puma_error.log
stderr_logfile_maxbytes=0
priority=20

[program:cron]
command=/usr/local/bin/supercronic /var/www/catima/config/crontab
user=catima
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/var/www/catima
environment=HOME="/home/catima",USER="catima"
priority=30
